# roles/github-worktree-cleanup/templates/github-worktree-cleanup.timer.j2
[Unit]
Description=GitHub Worktree Cleanup Timer
Requires=github-worktree-cleanup.service

[Timer]
OnCalendar={{ github_worktree_cleanup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
