# roles/github-worktree-cleanup/templates/github-worktree-cleanup.service.j2
[Unit]
Description=GitHub Worktree Cleanup Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecCondition=curl -sf --max-time 5 https://api.github.com >/dev/null 2>&1
ExecStart=%h/.local/bin/github-worktree-cleanup.sh
