# roles/github-worktree-cleanup/templates/github-worktree-cleanup.service.j2
[Unit]
Description=GitHub Worktree Cleanup Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/github-worktree-cleanup.sh
